# deploy/systemd/newblog-weread-sync.timer
[Unit]
Description=Run the newblog WeRead sync and reading briefing at 18:00 CST

[Timer]
OnCalendar=*-*-* 18:00:00 Asia/Shanghai
Persistent=true
AccuracySec=1s
Unit=newblog-weread-sync.service

[Install]
WantedBy=timers.target
