# deploy/systemd/newblog-weread-sync.service
[Unit]
Description=Synchronize the newblog bookshelf from WeRead
Wants=docker.service network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/workspace/newblog
ExecStart=/root/workspace/newblog/deploy/sync-weread.sh
TimeoutStartSec=35min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

