# deploy/systemd/newblog-weread-sync.timer
[Unit]
Description=Run the newblog WeRead bookshelf sync every day

[Timer]
OnCalendar=*-*-* 04:15:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10min
AccuracySec=1min
Unit=newblog-weread-sync.service

[Install]
WantedBy=timers.target
