# deploy/systemd/newblog-weread-sync.service
[Unit]
Description=Synchronize WeRead and send the daily reading briefing
Wants=docker.service network-online.target
After=docker.service network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/workspace/newblog
EnvironmentFile=/etc/newblog-reading-briefing.env
ExecStart=/usr/bin/python3 /root/workspace/newblog/scripts/reading-briefing.py sync-report
TimeoutStartSec=35min
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
